# pihole-backup.timer
[Unit]
Description=Daily backup of Pi-hole's configuration and lists

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
