# rclone-backup.timer
[Unit]
Description=Daily backup using rclone-backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
